# config/etc/udev/rules.d/50-sound.rules
ACTION=="add", SUBSYSTEM=="sound", RUN+="/usr/local/sbin/setup-sound"
